# run.crontab
*/5 * * * * /bin/sh /diy/ifttt-ewelink/run.sh
